# misc/tts.service
[Unit]
Description=TTS Service

[Service]
ExecStart=$PREFIX/bin/tts-service
KillMode=process
Restart=always

[Install]
WantedBy=default.target
