# misc/tts.service
[Unit]
Description=TTS Service

[Service]
ExecStart=/usr/local/bin/tts-service
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
